# cron.ini
[program:crond]
process_name=crond
command=crond
autostart=true
autorestart=true