# php-fpm.ini
[program:php-fpm]
process_name=php-fpm
command=php-fpm
autostart=true
autorestart=true